# debian/postinst
#!/bin/bash
set -e
#update-rc.d certmaster defaults > /dev/null
#invoke-rc.d certmaster start
#DEBHELPER#
