# debian/postinst
#!/bin/bash
set -e
update-rc.d certmasterd defaults > /dev/null
invoke-rc.d certmasterd start
#DEBHELPER#
