# debian/prerm
#!/bin/bash
set -e

PID=`ps -C "certmaster" --no-heading -o '%P %p'|awk '$1~/^1$/ {print$2;exit}'`
test -z "${PID}" || invoke-rc.d certmaster stop

#DEBHELPER#
